# agent-shell-base/etc/services.d/sshd/finish
#!/command/with-contenv bash
echo "[s6] sshd exited code=$1 (signal=$2)"
exit 0
