# agent-shell-base/etc/services.d/sshd/run
#!/command/with-contenv bash
exec /usr/sbin/sshd -f /opt/sshd_config -D -e
