# agent-shell-base/etc/services.d/supercronic/finish
#!/usr/bin/with-contenv bash
echo "[s6] supercronic exited code=$1 (signal=$2)"
exit 0
