# agent-shell-base/etc/services.d/supercronic/run
#!/command/with-contenv bash
exec supercronic "${AGENT_HOME}/.crontab"
